# rpi-crashplan/crashplan.conf
[program:crashplan]
directory=/opt/pi-config
command=/opt/pi-config/start_crashplan.sh
redirect_stderr=true
numprocs=1
autostart=true
autorestart=true
